# systemd/.config/systemd/user/dropbox.service
[Unit]
Description=Dropbox

[Service]
Type=simple
ExecStart=/home/bonext/.dropbox-dist/dropboxd

[Install]
WantedBy=default.target